# devices/xunison_exigo_hub_d50_5g/files/etc/init.d/modem-watchdog
#!/bin/sh /etc/rc.common
START=98

start() {
    /usr/bin/modem-watchdog.sh &
}
